# experiments/sleep_cycle/host/light_sleep.service
[Unit]
Description=Startup Script
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=olin
WorkingDirectory=/home/olin
ExecStart=/home/olin/VLA_Star/experiments/light_sleep.sh

[Install]
WantedBy=multi-user.target
