# experiments/light_sleep.service
[Unit]
Description=Startup Script
After=network.target

[Service]
Type=simple
ExecStart=/home/olin/VLA_Star/experiments/light_sleep.sh
Restart=no

[Install]
WantedBy=multi-user.target
